# ton618.service
[Unit]
Description=TON-618 Knowledge Singularity
After=network.target ollama.service

[Service]
Type=simple
User=giobon
WorkingDirectory=/home/giobon/CODE/ton618
ExecStart=/home/giobon/CODE/ton618/vortex_bin
Restart=always
RestartSec=10

# Variáveis de Ambiente (ajustadas para seu sistema DietPi)
Environment=PORT=6180
Environment=DOCS_DIR=/home/giobon/CODE/ton618/docs
Environment=WEB_DIR=/home/giobon/CODE/ton618/web/dist
Environment=BLEVE_INDEX_DIR=/home/giobon/ton618_data/ton618.bleve
Environment=STATE_DIR=/home/giobon/ton618_data
Environment=AUTH_USER=admin
Environment=AUTH_PASS=ton618_secret
Environment=OLLAMA_HOST=http://localhost:11434
Environment=OLLAMA_MODEL=nomic-embed-text

[Install]
WantedBy=multi-user.target
